# templates/systemd/oreo-metrics.service
[Unit]
Description=Oreo Cloud metrics collector

[Service]
Type=oneshot
User=oreo
Group=oreo
WorkingDirectory=/srv/oreo-cloud
ExecStart=/usr/bin/python3 /srv/oreo-cloud/control-plane/monitoring/collect_metrics.py
